# templates/neo4j_supervisor.conf.j2
[program:neo4j]
command=/usr/local/neo4j/bin/neo4j console
startsecs=60
numprocs=1
autostart=true
autorestart=true
